# roles/zapret/templates/zapret-watchdog.service.j2
[Unit]
Description=Zapret watchdog check

[Service]
Type=oneshot
ExecStart={{ zapret_dir }}/watchdog.sh
